# shell/.config/systemd/user/1password-ssh-agent.service
[Unit]
Description=1Password SSH Agent Bridge
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/1password-ssh-agent-bridge
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
